# units/etc/sv/mi_api/run
#!/bin/bash
set -m

exec /usr/bin/python3 /opt/mi_api.py